# genesis-blex/supervisor/conf.d/blockexplorer.conf
[program:blockexplorer]
command = /start_blockexplorer.sh
directory = /genesis-blex
stdout_logfile = /var/log/blockexplorer-stdout.log
stderr_logfile = /var/log/blockexplorer-stderr.log
user = root
