# genesis-blex/supervisor/conf.d/blockexplorer-worker.conf
[program:blockexplorer-worker]
command = /path/to/venv/bin/celery -B -A genesis_block_explorer.celery.tasks worker
directory = /path/to/project
user = nobody
