# genesis-blex/supervisor/conf.d/blockexplorer.conf
[program:blockexplorer]
command = /path/to/venv/bin/gunicorn -b 0.0.0.0:8000 --worker-class eventlet -w 1 genesis_block_explorer.socketio:app
directory = /path/to/project
user = nobody
