# genesis-bf/supervisor/conf.d/nginx.conf
[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;' -c /etc/nginx/nginx.conf
user = root
stdout_events_enabled = true
stderr_events_enabled = true
autorestart = true
